# barcodebuddy/root/etc/services.d/1websocketserver/run
#!/usr/bin/with-contenv bash

exec sudo -u root /usr/bin/php /app/bbuddy/wsserver.php
